# com.bioratio.paperfinder.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bioratio.paperfinder.sync</string>

    <!-- venv python + the sync script (absolute paths) -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jchan/Projects/PaperFinder/paper-finder/venv/bin/python3</string>
        <string>/Users/jchan/Projects/PaperFinder/paper-finder/examples/drive_sync.py</string>
    </array>

    <!-- run from the repo root so .env, service_account.json, and the DB resolve -->
    <key>WorkingDirectory</key>
    <string>/Users/jchan/Projects/PaperFinder/paper-finder</string>

    <!-- nightly at 02:30; if asleep then, launchd runs it on the next wake -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>2</integer>
        <key>Minute</key><integer>30</integer>
    </dict>

    <!-- keep it out of your way: background CPU/IO throttling -->
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>5</integer>

    <!-- don't run the moment it's loaded; scheduled + manual kickstart only -->
    <key>RunAtLoad</key>
    <false/>

    <!-- one append-only log for both stdout and stderr -->
    <key>StandardOutPath</key>
    <string>/Users/jchan/Library/Logs/paperfinder-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jchan/Library/Logs/paperfinder-sync.log</string>
</dict>
</plist>
